# skygrid-worker/supervisor.cfg
[program:skygrid_docker_worker]
command = /opt/recast-cloudconfig/skygrid-worker/skygrid-docker-worker
environment = DOCKER_WORKER_CONFIG="/opt/recast-cloudconfig/skygrid-worker/worker.cfg"
user = robot-worker
stderr_logfile = /var/log/skygrid-docker-worker-stderr.log
stdout_logfile = /var/log/skygrid-docker-worker-stdout.log
stopsignal = QUIT
