# skygrid-metascheduler/supervisor.cfg
[program:metascheduler]
command = /opt/recast-cloudconfig/skygrid-metascheduler/metascheduler
environment = METASCHEDULER_CONFIG="/opt/recast-cloudconfig/skygrid-metascheduler/metascheduler.cfg"
user = robot-metascheduler
stderr_logfile = /var/log/metascheduler-stderr.log
stdout_logfile = /var/log/metascheduler-stdout.log